# pi-gen-stage/00-signaeos-install/files/etc/systemd/system/companion-satellite.service
[Unit]
Description=Bitfocus Companion Satellite
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/companion-satellite
ExecStart=/opt/companion-satellite/companion-satellite
Restart=always
RestartSec=10
Environment=COMPANION_SATELLITE_CONFIG=/data/signaeos/satellite.json
StandardOutput=journal
StandardError=journal
SyslogIdentifier=companion-satellite

[Install]
WantedBy=multi-user.target
