# pi-gen-stage/00-signaeos-install/files/etc/systemd/system/signaeos-update.service
[Unit]
Description=SignageOS Update Check

[Service]
Type=oneshot
ExecStart=/usr/bin/signaeos-update check
StandardOutput=journal
