# rootfs/etc/systemd/system/signaeos-display1.service
[Unit]
Description=SignageOS Display 1 (Web)
After=network.target signaeos-webui.service
Wants=signaeos-webui.service

[Service]
Type=simple
User=admin
ExecStart=/usr/bin/signaeos-display1
Restart=always
RestartSec=5
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-1
Environment=HOME=/home/admin
StandardOutput=journal
StandardError=journal
SyslogIdentifier=signaeos-display1

[Install]
WantedBy=multi-user.target
