# rootfs/etc/systemd/system/signaeos-display1.service
[Unit]
Description=SignageOS Display 1 (Web)
After=network.target sway.service signaeos-webui.service
Wants=sway.service signaeos-webui.service

[Service]
Type=simple
User=signaeos
Group=signaeos
ExecStartPre=+/usr/bin/install -d -o signaeos -g signaeos -m 0755 /run/signaeos
ExecStart=/usr/bin/signaeos-display1
Restart=always
RestartSec=5
Environment=XDG_RUNTIME_DIR=/run/signaeos-runtime
Environment=HOME=/home/signaeos
StandardOutput=journal
StandardError=journal
SyslogIdentifier=signaeos-display1

[Install]
WantedBy=multi-user.target
