# rootfs/etc/systemd/system/signaeos-display2.service
[Unit]
Description=SignageOS Display 2 (NDI)
After=signaeos-display1.service
Wants=signaeos-display1.service

[Service]
Type=simple
User=admin
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/signaeos-display2
Restart=always
RestartSec=10
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-1
Environment=HOME=/home/admin
StandardOutput=journal
StandardError=journal
SyslogIdentifier=signaeos-display2

[Install]
WantedBy=multi-user.target
